# cloudflared-launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.avimarketdata.cloudflared</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/avi.horowitz/Documents/LayoffTracker -74 - with earning event copy/bin/cloudflared</string>
    <string>tunnel</string>
    <string>run</string>
    <string>avi-marketdata-tunnel</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/avi.horowitz/Documents/LayoffTracker -74 - with earning event copy</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/avi.horowitz/Documents/LayoffTracker -74 - with earning event copy/cloudflared.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/avi.horowitz/Documents/LayoffTracker -74 - with earning event copy/cloudflared.log</string>
</dict>
</plist>
